# hm-diag.service
[Unit]
Description=start diagnostic service

[Service]
ExecStart=/bin/bash /opt/diag_start.sh

[Install]
WantedBy=multi-user.target